# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/data/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/data/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///data/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mysql]
priority=10
command=/usr/sbin/mysqld --datadir=/data/mysql --socket=/data/run/mysqld/mysqld.sock --pid-file=/data/run/mysqld/mysqld.pid --port=3306 --bind-address=127.0.0.1 --skip-name-resolve --mysqlx=0 --log-error=/data/logs/mysql-error.log
directory=/data
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/data/logs/mysql.log
stderr_logfile=/data/logs/mysql.err

[program:redis]
priority=20
command=/usr/bin/redis-server /data/run/redis.conf
directory=/data
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.err

[program:nocodb]
priority=40
command=/usr/local/bin/nocodb
directory=/data/nocodb
environment=HOME="/home/user"
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/data/logs/nocodb.log
stderr_logfile=/data/logs/nocodb.err

[program:ops-service]
priority=50
command=/usr/bin/python3 /usr/local/bin/db-ops-service
directory=/home/user
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
priority=80
command=/usr/sbin/nginx -g "daemon off;"
directory=/home/user
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
